# deploy/naviproxy.service
[Unit]
Description=NaviProxy dashboard and reverse proxy controller
After=network-online.target caddy.service
Wants=network-online.target caddy.service

[Service]
Type=simple
WorkingDirectory=/opt/naviproxy
EnvironmentFile=/etc/naviproxy/naviproxy.env
ExecStart=/opt/naviproxy/start.sh
Restart=on-failure
RestartSec=5
User=naviproxy
Group=naviproxy

[Install]
WantedBy=multi-user.target
